# scripts/guardian-one.service
# Guardian One — systemd service unit
#
# Install:
#   # 1. Set your project path and generate the service file
#   sed 's|/GUARDIAN_DIR|/home/youruser/JT|g' scripts/guardian-one.service > /tmp/guardian-one.service
#
#   # 2. Create dedicated service user (or change User/Group to your own)
#   sudo useradd -r -s /usr/sbin/nologin guardian
#
#   # 3. Ensure startup script is executable
#   chmod +x scripts/guardian_daemon.sh
#
#   # 4. Ensure the guardian user can access the project directory
#   #    Option A: Run as your own user — change User/Group in the service file
#   #    Option B: Grant guardian access — sudo chown -R guardian:guardian ~/JT
#
#   # 5. Install and enable
#   sudo cp /tmp/guardian-one.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable guardian-one    # start on boot
#   sudo systemctl start guardian-one     # start now
#
#   # 6. Verify it is running
#   sudo systemctl status guardian-one
#   sudo journalctl -u guardian-one --no-pager -n 20
#
# Manage:
#   sudo systemctl status guardian-one    # check health
#   sudo journalctl -u guardian-one -f    # tail logs
#   sudo systemctl restart guardian-one   # restart
#   sudo systemctl stop guardian-one      # graceful stop (SIGTERM)
#
# Paths:
#   All paths use /GUARDIAN_DIR as a placeholder. Replace with your actual
#   project path using the sed command in step 1 above.
#
# Security:
#   - Runs as dedicated 'guardian' user (or change User/Group to your own)
#   - Private /tmp, no privilege escalation, restricted filesystem
#   - Credentials loaded from EnvironmentFile (.env), never in this unit

[Unit]
Description=Guardian One — Autonomous Life Management System
Documentation=https://github.com/CuriousCat777/JT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=guardian
Group=guardian
WorkingDirectory=/GUARDIAN_DIR
EnvironmentFile=-/GUARDIAN_DIR/.env
ExecStart=/GUARDIAN_DIR/scripts/guardian_daemon.sh

# ── Restart policy ──────────────────────────────────────────────────
# Auto-restart on crash with a fixed 5-second delay between attempts
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# ── Graceful shutdown ───────────────────────────────────────────────
# Send SIGTERM first, wait 30s for clean state persistence, then SIGKILL
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# ── Environment ─────────────────────────────────────────────────────
# Secrets and API keys are loaded by the startup script from .env
# but can also be set here for systemd-native credential management
EnvironmentFile=-/GUARDIAN_DIR/.env

# ── Logging ─────────────────────────────────────────────────────────
# stdout/stderr → journald (view with: journalctl -u guardian-one)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=guardian-one

# ── Security hardening ──────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/GUARDIAN_DIR/data /GUARDIAN_DIR/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
LockPersonality=true

[Install]
WantedBy=multi-user.target
